# pr-status.service
[Unit]
Description=Check GitHub PR status

[Service]
Type=simple
ExecStart=bash /usr/local/bin/pr-status.sh

[Install]
WantedBy=default.target
